# autostart/ignore-lid-switch-tweak.desktop
[Desktop Entry]
Type=Application
Name=ignore-lid-switch-tweak
Exec=/usr/libexec/gnome-tweak-tool-lid-inhibitor

